# real_estate_dashboard/backend/scripts/systemd/market-intelligence-update.service
[Unit]
Description=Market Intelligence Daily Update
After=network.target postgresql.service redis.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/home/user/real_estate_dashboard/backend
Environment="PYTHONPATH=/home/user/real_estate_dashboard/backend"
ExecStart=/usr/local/bin/python3 /home/user/real_estate_dashboard/backend/scripts/update_market_intelligence.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=market-intelligence-update

# Restart on failure
Restart=on-failure
RestartSec=300

# Timeout after 30 minutes
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
